# crontab
15 1 * * 6 root /usr/local/bin/python /main.py
# an empty line
