# crontab
15 1 * * * root /usr/local/bin/python /main.py
# an empty line
